# agent-service/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:neo4j]
command=/var/lib/neo4j/bin/neo4j console
autostart=true
autorestart=true
startsecs=30
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/neo4j.log
stderr_logfile=/var/log/supervisor/neo4j-error.log
environment=NEO4J_HOME="/var/lib/neo4j",NEO4J_DATA="/data/neo4j"
priority=10

[program:agent-service]
command=uvicorn agent_service.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/agent-service.log
stderr_logfile=/var/log/supervisor/agent-service-error.log
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app",NEO4J_URI="bolt://localhost:7687"
priority=20
startretries=5
